# biglinux-apps-rename/usr/lib/systemd/user/fonografo.service
[Unit]
Description=fonografo

[Service]
Type=simple
#StandardOutput=tty
ExecStart=/usr/bin/fonografo

[Install]
WantedBy=multi-user.target 
