# biglinux-apps-rename/usr/lib/systemd/user/fonografo.timer
 
[Unit]
Description=fonografo

[Timer]
OnUnitActiveSec=60min

[Install]
WantedBy=multi-user.target 
